# backend_esempio/estrattore.service
# Esempio di servizio systemd per tenere il backend sempre attivo.
# Copiare in /etc/systemd/system/estrattore.service, adattare i percorsi, poi:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now estrattore
#   sudo systemctl status estrattore

[Unit]
Description=Estrattore testi sito - backend FastAPI (LAN)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# utente non-root dedicato: sudo useradd -r -s /usr/sbin/nologin -G intrawelt estrattore
# (il repository resta sotto la home di intrawelt, quindi estrattore va aggiunto al gruppo
# "intrawelt" per ereditare i permessi di gruppo gia' presenti su home/Scrivania/repo).
User=estrattore
WorkingDirectory=/home/intrawelt/Scrivania/website-analyst
Environment=OUTPUT_BASE=/srv/output
Environment=ARCHIVE_BASE=/mnt/downloaded-websites
# senza questa, l'utente estrattore non troverebbe il Chromium di Playwright,
# scaricato nella cache dell'utente intrawelt (serve anche "chmod g+rx" su
# /home/intrawelt/.cache, vedi deployment.md)
Environment=PLAYWRIGHT_BROWSERS_PATH=/home/intrawelt/.cache/ms-playwright
# usa il python del virtualenv
ExecStart=/home/intrawelt/Scrivania/website-analyst/.venv/bin/uvicorn backend_esempio.app:app --host 192.168.20.24 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
